# shared/ovldaemon.conf
#
# Upstart script for ina-server
#

description "The VLF observatory network daemon"

exec /usr/bin/ovldaemon

start on runlevel [23]
stop on runlevel [016]

